# cmd/looper/rc.local
#!/bin/sh -e
# /etc/rc.local
# This script launches audio looper at boot time
printf "rc.local started"      # show start of execution
set -x                         # tell sh to display commands before execution

PLAYERPATH=/home/pi/go/src/github.com/ausocean/av/cmd/looper
cd $PLAYERPATH

sudo ./looper -path=/home/pi/example.mp3 &

exit 0
